# deploy/nounours-scoreboard.service
[Unit]
Description=nounours.app — serveur de classement WebDeGentil (isole, hors moteur)
After=network.target

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/nounours-scoreboard
# Panneau admin (optionnel) : définir la clé secrète — JAMAIS commitée — via :
#   sudo systemctl edit nounours-scoreboard  →  [Service]  Environment=NOUNOURS_ADMIN_KEY=...
# Sans cette variable, toutes les routes /admin/* renvoient 403 (admin désactivé).
ExecStart=/usr/bin/python3 /home/debian/nounours-scoreboard/server.py 8790
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
